# meta-riscv/recipes-phosphor/console/obmc-console/obmc-console-ttyS2.socket
[Unit]
Description= Console SSH socket for Slot2
Wants=obmc-console@ttyS2.service

[Socket]
ListenStream=2302
Accept=yes

[Install]
WantedBy=sockets.target
